# com.webforge.orquestador.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.webforge.orquestador</string>

  <!-- Ejecutar a las 08:00 todos los días -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/davidnicolassoto/webforge/orquestador/run-daily.sh</string>
  </array>

  <!-- Logs (stdout y stderr van al mismo fichero) -->
  <key>StandardOutPath</key>
  <string>/Users/davidnicolassoto/Library/Logs/WebForge/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/davidnicolassoto/Library/Logs/WebForge/launchd.log</string>

  <!-- Reintentar si el Mac estaba apagado a las 08:00 -->
  <key>RunAtLoad</key>
  <false/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/davidnicolassoto</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
